# deploy/stock-trader.service
[Unit]
Description=KIS Auto Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rocky
WorkingDirectory=/home/rocky/stock-auto-trader
Environment=PATH=/home/rocky/stock-auto-trader/.venv/bin:/usr/bin
ExecStart=/home/rocky/stock-auto-trader/.venv/bin/python main.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
